# service/mci-analyzer.service
[Unit]
Description=Multiple CI - Analyzer
Wants=network.target

[Service]
Type=simple
User=root
ExecStart=/bin/bash -c "source /root/multiple-ci/venv/bin/activate && mci-analyzer"
Restart=on-failure

[Install]
WantedBy=multi-user.target
